# supervisor/configs-supervisor/autobet armando.conf
[program:autobet_armando]
process_name=%(program_name)s_%(process_num)02d
command=node /autobet/src/estrategias/escanteios_funil_armando.js
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/autobet_armando.log